# fetch_joke.service
[Unit]
Description=Fetch Chuck Norris Joke

[Service]
Type=oneshot
ExecStart=/usr/bin/python3 /usr/local/bin/fetch_joke.py

